# udev.rules
SUBSYSTEMS=="input", ATTRS{name}=="Nintendo Wii Remote Balance Board", RUN+="/bin/systemctl start wiiscale.service"